# configs/etc/cron.daily/backup
#!/bin/bash

SCRIPT=/opt/backup.sh
RUN_TYPE=scheduled

if [[ -f $SCRIPT ]]; then
  source "$SCRIPT"
fi
